# tmp/cron/crontabs/root
* * * * * /usr/bin/php-cgi -q /mnt/plugin/crpt/mk_ct_report.cgi
